# infrastructure/ansible/roles/seedbox.os/templates/cross-seed.service.j2
[Unit]
Description=cross-seed

[Service]
TimeoutStartSec=0
User=qbit
ExecStart=/usr/bin/cross-seed daemon

[Install]
WantedBy=multi-user.target

[Service]
Restart=on-failure
RestartSec=5s
