# infrastructure/ansible_old/nas/playbooks/templates/zfs_exporter.service.j2
[Unit]
Description=ZFS Exporter
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/zfs_exporter
Restart=always
User=durandal
Group=durandal

[Install]
WantedBy=multi-user.target
